# systemd/.config/systemd/user/eye-break.service
[Unit]
Description=Eye Break Reminder

[Service]
Type=oneshot
ExecStart=/home/jho/.local/bin/eye-break.sh
