# systemd/.config/systemd/user/eye-break.timer
[Unit]
Description=Eye Break Timer

[Timer]
OnBootSec=5m
OnUnitActiveSec=2h
Unit=eye-break.service

[Install]
WantedBy=timers.target
