# systemd/.config/systemd/user/random-wallpaper.service
[Unit]
Description=Change wallpaper randomly

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/random-wallpaper.sh
