# systemd/.config/systemd/user/standup-reminder.service
[Unit]
Description=Stand Up Reminder

[Service]
Type=oneshot
ExecStart=/home/jho/.local/bin/standup-reminder.sh
