# systemd/.config/systemd/user/water-reminder.service
[Unit]
Description=Water Reminder

[Service]
Type=oneshot
ExecStart=/home/jho/.local/bin/water-reminder.sh
