# systemd/.config/systemd/user/water-reminder.timer
[Unit]
Description=Water Reminder Timer

[Timer]
OnBootSec=15m
OnUnitActiveSec=45m
Unit=water-reminder.service

[Install]
WantedBy=timers.target
